# FichiersIntallation/rc.local
sudo /home/user/NirScanSweep/Main.py 

exit 0 